# templates/etc/service/getty-TTY/finish
#!/bin/sh
# If agetty exits, sleep briefly to avoid busy-loop
exec sleep 5
